# supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:fastapi]
command=python -m uvicorn api.main:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
stderr_logfile=/tmp/fastapi.err.log
stdout_logfile=/tmp/fastapi.out.log

[program:streamlit]
command=streamlit run dashboard/dashboard.py --server.port 7860 --server.address 0.0.0.0 --server.headless true
directory=/app
autostart=true
autorestart=true
stderr_logfile=/tmp/streamlit.err.log
stdout_logfile=/tmp/streamlit.out.log